# docker/supervisor/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D
autorestart=true

[program:mysql]
command=/usr/local/bin/docker-entrypoint.sh mysqld
